# scripts/day.timer
[Unit]
Description=Day/Night detect

[Timer]
OnBootSec=0min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target